# devopsfetch.service
[Unit]
Description=DevOpsFetch Service
After=network.target

[Service]
ExecStart=/path/to/devopsfetch.sh
Restart=always

[Install]
WantedBy=multi-user.target
