# devopsfetch.service
[Unit]
Description=DevOpsFetch Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /path/to/devopsfetch.py -l
Restart=always

[Install]
WantedBy=multi-user.target
